# hdmi_sound_switch.rules
SUBSYSTEM=="drm", RUN+="/usr/lib/hdmi_sound_switch/hdmi_sound_switch"